# pwbuild/ap/clamav/clamav.cron
#!/bin/sh

/usr/bin/freshclam --datadir /usr/share/clamav
